# packaging/macos/com.telemon.exporter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.telemon.exporter</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/libexec/telemon/telemon-exporter</string>
      <string>run</string>
      <string>--config</string>
      <string>/Library/Application Support/Telemon/exporter.yml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Library/Logs/Telemon/exporter.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Logs/Telemon/exporter.err.log</string>
  </dict>
</plist>
